# scripts/liveoilprices.service
[Unit]
Description=Live Oil Prices Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/liveoilprices
ExecStart=/home/deploy/liveoilprices/bin/server
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=liveoilprices

Environment=PORT=8080

# Optional secrets file. Contains lines like `EIA_API_KEY=xxxxx`.
# Create it on the server with:
#   sudo install -d -m 700 -o deploy -g deploy /etc/liveoilprices
#   sudo tee /etc/liveoilprices/env >/dev/null <<<'EIA_API_KEY=your-key'
#   sudo chmod 600 /etc/liveoilprices/env
# The leading "-" means systemd won't fail if the file is absent.
EnvironmentFile=-/etc/liveoilprices/env

NoNewPrivileges=true
ProtectHome=no
PrivateTmp=true

LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
