# scripts/liveoilprices.service
[Unit]
Description=Live Oil Prices Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/liveoilprices
ExecStart=/home/deploy/liveoilprices/bin/server
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=liveoilprices

Environment=PORT=8080

NoNewPrivileges=true
ProtectHome=no
PrivateTmp=true

LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
